# systemd/time-sync.service
[Unit]
Description=Sync Clockify time entries into Toggl Track
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/time-sync/time-sync.env
ExecStart=/opt/time-sync/.venv/bin/time-sync
User=time-sync
Group=time-sync
StateDirectory=time-sync
Environment=TIME_SYNC_STATE_FILE=/var/lib/time-sync/state.json
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
